# supervisor-gifbox.conf
[program:gifbox]
command=/home/tim/gifbox/gunicorn.sh
directory=/home/tim
user=tim
autostart=true
autorestart=true
stdout_logfile = /var/log/supervisor/gifbox.log
redirect_stderr=true
