# supervisor-djropbox.conf
[program:djropbox]
command=/home/tim/djropbox/gunicorn.sh
directory=/home/tim
user=tim
autostart=true
autorestart=true
stdout_logfile = /var/log/supervisor/djropbox.log
redirect_stderr=true